# workflows/bootc-debs-builder/preinstall/output/bootupd/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
  configure | triggered)
    deb-systemd-helper enable bootloader-update.service

    # Assemble the EFI components bootupd installs to the ESP. grub is staged by
    # the grub-efi-amd64-signed package under /usr/lib/efi/grub2/...; here we add
    # shim (from shim-signed, a Pre-Depend) under /usr/lib/efi/shim/<ver>/EFI in
    # the 3-level <name>/<version>/EFI layout bootupd reads (cf. bootupd v0.2.35
    # src/efi.rs and workflows/image-tests/tests.yml), then generate the update
    # metadata bootupd needs to actually install the bootloader. This runs at
    # image build time (writable /usr); bootc systems never apt-install on a
    # booted, read-only /usr.
    shim_ver=$(dpkg-query -W -f='${Version}' shim-signed 2>/dev/null || echo 0)
    efidir="/usr/lib/efi/shim/${shim_ver}/EFI"
    mkdir -p "${efidir}/BOOT" "${efidir}/debian"
    # shimx64.efi.signed is the boot entrypoint: it lands as the removable-media
    # fallback EFI/BOOT/BOOTX64.EFI and as the vendor EFI/debian/shimx64.efi.
    cp /usr/lib/shim/shimx64.efi.signed "${efidir}/BOOT/BOOTX64.EFI"
    cp /usr/lib/shim/shimx64.efi.signed "${efidir}/debian/shimx64.efi"
    # MokManager, for MOK enrollment; best-effort. Only launched on a pending
    # MOK request, so harmless. We deliberately do NOT ship fbx64.efi (the shim
    # fallback): booted from the removable path with no NVRAM entry, shim would
    # launch it, and with no BOOTX64.CSV to rebuild the entry from it just resets
    # the machine in an endless loop. Without fbx64, shim loads grubx64.efi from
    # its own directory directly, which is exactly what we want. bootupd already
    # owns the real NVRAM entry via efibootmgr, so the fallback is redundant.
    cp /usr/lib/shim/mmx64.efi.signed "${efidir}/BOOT/mmx64.efi" 2>/dev/null || true

    # bootupd v0.2.35 write_update_metadata does open_dir(usr/lib/bootupd/updates)
    # (not create) and fails ENOENT if it is missing; the deb ships no empty dir.
    mkdir -p /usr/lib/bootupd/updates

    bootupctl backend generate-update-metadata /

    # generate-update-metadata re-stages the shim fallback fbx64.efi but not the
    # BOOTX64.CSV it needs. Booted from the removable path \EFI\BOOT\BOOTX64.EFI
    # with no NVRAM entry, shim runs fbx64, and without a CSV telling it which
    # entry to (re)create it just resets the machine forever. Debian ships the
    # CSV at /usr/lib/shim/BOOTX64.CSV; stage it in the vendor dir so bootupd lays
    # it on the ESP and the fallback works as designed. After generate-update-
    # metadata so it is not clobbered.
    cp /usr/lib/shim/BOOTX64.CSV "${efidir}/debian/BOOTX64.CSV" \
      2>/dev/null || true
    ;;
esac
