# workflows/bootc-debs-builder/preinstall/output/ostree/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
  configure | triggered)
    deb-systemd-helper enable ostree-remount.service
    ;;
esac
