# workflows/bootc-debs-builder/firstbootpostinstall/usr/lib/systemd/system/getty@tty1.service.d/firstboot.conf
[Service]
ExecStartPre=-/usr/local/bin/firstboot-user-setup