# workflows/bootc-debs-builder/firstbootpostinstall/usr/lib/systemd/system/getty@tty1.service.d/firstboot.conf
[Service]
# The wizard runs as ExecStartPre and is paced by a human answering prompts.
# Without this, systemd's default 90s start timeout (TimeoutStartSec) kills it
# mid-setup, and getty's Restart=always then relaunches it from the very first
# question. Let it take as long as the operator needs.
TimeoutStartSec=infinity
ExecStartPre=-/usr/local/bin/firstboot-user-setup