# workflows/bootc-debs-builder/preinstall/output/firstboot-user-setup/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure|triggered)
        deb-systemd-helper enable firstboot-user-setup.service
        ;;
esac
