# workflows/bootc-debs-builder/ifupdown2postinstall/usr/lib/systemd/system/networking.service.d/ifupdown2-autoconf.conf
[Service]
ExecStartPre=/usr/sbin/ifupdown2-autoconf
ExecStartPost=/usr/sbin/domain-set